# files/etc/periodic/hourly/10-curator
#!/bin/sh

set -e

curator \
    --config "/etc/curator/config.yml" \
    "/etc/curator/actions-delete-old-indices.yml"
